# debian/docsis-pnm/lib/systemd/system/docsis-pnm.service
[Unit]
Description=DOCSIS PNM

[Service]
Type=simple
ExecStart=/usr/bin/docsis-pnm
User=docsis-pnm

[Install]
WantedBy=multi-user.target
